# check-expire.service
# /etc/systemd/system/check-expire.service
cat > /etc/systemd/system/check-expire.service << 'EOF'
[Unit]
Description=Check and disable expired users
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/check_users_expire.py

[Install]
WantedBy=multi-user.target
EOF
