# check-expire.service
cat > /etc/systemd/system/check-expire.service << 'EOF' 
[Unit]
Description=Check and disable expired users

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_users_expire.py
