# check-user-usage.timer
# /etc/systemd/system/check-user-usage.timer
cat > /etc/systemd/system/check-user-usage.timer << 'EOF'
[Unit]
Description=Check SSH User Traffic Hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=check-user-usage.service

[Install]
WantedBy=timers.target
EOF
