# check-expire.timer
[Unit]
Description=Run expire checker every 1 hour

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
