# check-user-usage.service
# /etc/systemd/system/check-usage.service
cat > /etc/systemd/system/check-user-usage.service << 'EOF'
[Unit]
Description=Run check_user_usage.py script
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/check_user_usage.py
EOF
